# deploy/railway/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:xvfb]
priority=10
command=/usr/bin/Xvfb :99 -screen 0 1280x1024x24 -nolisten tcp
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:openclaw-gateway]
priority=20
directory=/app
user=node
command=/usr/local/bin/docker-entrypoint.sh node openclaw.mjs gateway --allow-unconfigured
environment=HOME="/home/node",PORT="18789",OPENCLAW_AUTH_MODE="token"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:roboclaws-interactive]
priority=30
directory=/opt/roboclaws
command=/opt/roboclaws/scripts/appliance-run-interactive.sh
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_ROBOCLAWS_HOME)s"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:snapshot-viewer]
priority=40
directory=/opt/roboclaws
command=/opt/roboclaws/.venv/bin/python scripts/view-snapshots.py --run-dir /data/runs/current --host 127.0.0.1 --port 8787
environment=HOME="%(ENV_ROBOCLAWS_HOME)s"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
priority=50
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
